# lib/init/router-cron.conf
# cron - regular background program processing daemon
#

description	"regular background program processing daemon"

start on runlevel [2345]
stop on runlevel [!2345]

#expect fork
respawn

script
	exec >>/tmp/upstart.log 2>&1
	set -x
	export ROUTER_USER=www
	test -s "/app/www/init_env" && . /app/www/init_env || true
	rm -f /tmp/god.*.sock
	exec su -lc 'god -c /app/www/config/cron.god -D' $ROUTER_USER
end script
