# lib/init/router-cron.conf
# cron - regular background program processing daemon
#

description	"regular background program processing daemon"

start on runlevel [2345]
stop on runlevel [!2345]

#expect fork
respawn
console output

script
	rm -f /tmp/god.*.sock
	exec su -lc 'god -c $STATS_APP_PATH/config/cron.god -D' root >> /tmp/upstart.log 2>&1
	#exec /usr/local/rvm/rubies/ruby-1.9.2-p180/bin/ruby /usr/local/rvm/gems/ruby-1.9.2-p180/bin/god -c /home/ishi/Projects/router-monitor/config/cron.god >> /tmp/upstart.log 2>&1
end script
